# apps/scraper/launchd/com.eat-thing.scraper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.eat-thing.scraper</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/keneselautusi/Documents/Code/PROJECTS/eat-thing/apps/scraper/dist/index.js</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/keneselautusi/Documents/Code/PROJECTS/eat-thing/apps/scraper</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/keneselautusi/Library/Logs/eat-thing-scraper.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/keneselautusi/Library/Logs/eat-thing-scraper.error.log</string>
</dict>
</plist>
